# deploy/supervisord.conf
; Flock OS production supervisor (FLO-246 Phase 6.1).
;
; The multi-process manager that preserves the FLO-121 scaled-socketio tier
; across restarts — the explicit ADR blocker #2 ("Procfile preserves the scaled
; tier across dashboard restarts"). Frappe Cloud Server plan dashboard restarts
; invoke the VM's init; supervisor here is the single source of truth for the
; bench topology, so a dashboard restart brings the WHOLE tier back, not just
; gunicorn.
;
; Process tree:
;   supervisord
;   ├── nginx            (front edge: TLS-term + sticky-L7 socketio LB on :9000
;   │                     + reverse proxy to gunicorn on :8080)
;   ├── gunicorn         (Frappe web — metered, light per hosting-quote)
;   ├── bench-worker     (queue workers — unmetered; rides the FLO-15 bulk path)
;   ├── bench-schedule   (Frappe scheduler — fires hooks.py scheduled events)
;   ├── socketio-tier    (scripts/dev/scale-socketio.sh start N --lb nginx:
;   │                     spawns N node socketio backends + the nginx sticky-L7
;   │                     LB; idempotent + self-recovering)
;   └── redis-adapter-ping (liveness probe for the dedicated adapter Redis —
;                          logs a warn if FLOCK_SIO_ADAPTER_REDIS is unset, so a
;                          silent stall like FLO-127 §2 is loud, not silent)
;
; Owned by root (supervisor must start as root to bind :8080); programs drop
; privileges to the `frappe` user. Logs go to /var/log/flock-os/ for a clean
; `docker logs` / cloud log-stream contract.

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemons=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/flock-os

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ----- nginx: the front edge (TLS termination happens at Cloudflare/Frappe
;        Cloud edge; this nginx does the sticky-L7 socketio LB + web proxy).
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/flock-os/nginx.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/flock-os/nginx.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
priority=10

; ----- gunicorn: Frappe web tier (bound to loopback; nginx proxies :8080 → :8000).
;        Worker count is tunable via GUNICORN_WORKERS at container start. The WSGI
;        entry is `frappe.app:application` (Frappe's own WSGI app — Frappe apps
;        don't ship a per-app app.py), matching the proven docker dev tier.
[program:gunicorn]
command=gunicorn \
    --bind 127.0.0.1:8000 \
    --workers %(ENV_GUNICORN_WORKERS)s \
    --threads %(ENV_GUNICORN_THREADS)s \
    --timeout %(ENV_GUNICORN_TIMEOUT)s \
    --preload \
    --max-requests %(ENV_GUNICORN_MAX_REQUESTS)s \
    --max-requests-jitter %(ENV_GUNICORN_MAX_REQUESTS_JITTER)s \
    frappe.app:application
directory=/home/frappe/frappe-bench
user=frappe
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=30
stdout_logfile=/var/log/flock-os/gunicorn.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/flock-os/gunicorn.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
priority=20

; ----- bench worker: background queue processing (FLO-15 bulk path rides here).
[program:bench-worker]
command=bench worker --queue short,long,default
directory=/home/frappe/frappe-bench
user=frappe
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/flock-os/worker.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/flock-os/worker.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
priority=30

; ----- bench scheduler: Frappe scheduled events (hooks.py).
[program:bench-schedule]
command=bench schedule
directory=/home/frappe/frappe-bench
user=frappe
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stopwaitsecs=30
stdout_logfile=/var/log/flock-os/schedule.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/flock-os/schedule.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
priority=40

; ----- socketio tier: scripts/dev/scale-socketio.sh start N --lb external.
;        Idempotent + self-recovering; brings up N node backends on :9001..900N.
;        --lb external means NO LB process is started here — the prod nginx
;        (the sibling supervisor program above) IS the sticky-L7 LB on :9000,
;        with its upstream generated from FLOCK_SIO_PROCESSES by entrypoint.sh.
;        This avoids two processes fighting over :9000 (architect fix FLO-246).
;        FLOCK_SIO_PROCESSES (default 4) sets the backend count;
;        FLOCK_SIO_ADAPTER_REDIS points at the dedicated adapter Redis (FLO-127
;        §2). The scale-socketio.sh process is itself the parent that supervises
;        the N backends; if it dies, supervisor restarts the whole tier cleanly
;        (stop-then-start is idempotent in the script).
[program:socketio-tier]
command=bash -lc 'flock-os-scale-socketio start "${FLOCK_SIO_PROCESSES:-4}" --lb external'
directory=/home/frappe/frappe-bench
user=frappe
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1",FLOCK_SIO_LB_PORT="9000",FLOCK_SIO_BASE_PORT="9001",FLOCK_SIO_BENCH="/home/frappe/frappe-bench"
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stopwaitsecs=15
stdout_logfile=/var/log/flock-os/socketio-tier.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/flock-os/socketio-tier.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
priority=50

; ----- redis-adapter-ping: liveness probe for the dedicated adapter Redis
;        (FLO-127 §2). The @socket.io/redis-adapter silently degrades to
;        single-node fan-out if it loses its adapter Redis; this probe makes that
;        stall LOUD: a real adapter Redis that stops answering PING exits non-zero
;        and supervisor marks it FATAL (visible in `supervisorctl status` /
;        docker logs), instead of a silent cross-node blackout.
;
;        - FLOCK_SIO_ADAPTER_REDIS unset → warn-only loop (single-node is a
;          valid launch-band shape; the probe stays RUNNING so it isn't noise).
;          Set FLOCK_SIO_ADAPTER_PING_DISABLE=1 to skip the probe entirely.
;        - FLOCK_SIO_ADAPTER_REDIS set but PING fails → ERROR + exit 1 (LOUD).
;        Interval is FLOCK_SIO_ADAPTER_PING_INTERVAL (default 15s). redis-cli
;        comes from the redis-tools apt package baked into the image.
[program:redis-adapter-ping]
command=bash -lc '\
    interval="${FLOCK_SIO_ADAPTER_PING_INTERVAL:-15}"; \
    if [ "${FLOCK_SIO_ADAPTER_PING_DISABLE:-0}" = "1" ]; then \
        echo "[redis-adapter-ping] disabled (FLOCK_SIO_ADAPTER_PING_DISABLE=1)"; exit 0; \
    fi; \
    while true; do \
        if [ -z "${FLOCK_SIO_ADAPTER_REDIS:-}" ]; then \
            echo "[redis-adapter-ping] WARN: FLOCK_SIO_ADAPTER_REDIS unset — clustered adapter Redis not wired; realtime fan-out is single-node only (FLO-127 §2)."; \
        else \
            if ! redis-cli -u "${FLOCK_SIO_ADAPTER_REDIS}" ping >/dev/null 2>&1; then \
                echo "[redis-adapter-ping] ERROR: adapter Redis PING failed (${FLOCK_SIO_ADAPTER_REDIS}) — socketio cross-node fan-out is stalled (FLO-127 §2)."; \
                exit 1; \
            fi; \
        fi; \
        sleep "${interval}"; \
    done'
directory=/home/frappe/frappe-bench
user=frappe
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
autostart=true
autorestart=true
startsecs=0
stopasgroup=true
killasgroup=true
stopwaitsecs=5
stdout_logfile=/var/log/flock-os/redis-adapter-ping.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/flock-os/redis-adapter-ping.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=60
